# dot_config/systemd/user/ssh-agent.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ssh-agent -D -t 12h
